# ops/pjsk-renderer.service
[Unit]
Description=PJSK Render Service (FastAPI + Playwright)
After=network.target

[Service]
Type=simple
User=pjsk
WorkingDirectory=/opt/pjsk-astrbot/current
Environment=RENDER_HOST=127.0.0.1
Environment=RENDER_PORT=3000
Environment=RENDER_MAX_CONCURRENT=4
ExecStart=/opt/pjsk-astrbot/current/.venv/bin/python -m uvicorn render_service.main:app --host 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=5

# Security
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
