# ops/pjsk-renderer.service
[Unit]
Description=PJSK Render Service (FastAPI + Playwright)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/data/plugins/astrbot_plugin_pjsk
Environment=RENDER_HOST=127.0.0.1
Environment=RENDER_PORT=3000
Environment=RENDER_MAX_CONCURRENT=4
ExecStart=/root/.local/share/uv/tools/astrbot/bin/python render_service/main.py
Restart=on-failure
RestartSec=5

# Security
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
